# playbook/templates/etc_systemd_system_mjpgstream.service
[Unit]
Description=mjpg streamer

[Service]
ExecStart=/usr/local/appscript/mjpg.sh

[Install]
WantedBy=multi-user.target
