# playbook/templates/etc_systemd_system_aria2.service
[Unit]
Description=aria2 rpc

[Service]
ExecStart=/usr/local/appscript/aria2.sh

[Install]
WantedBy=multi-user.target
